# graviton.service
[Unit]
Description=Graviton Server

[Service]
ExecStart=/usr/bin/graviton-server

[Install]
WantedBy=default.target
